# builder/fcnet
#!/sbin/openrc-run

command="/bin/fcnet-setup"
